# backup/crontab
# Запускаем резервное копирование в 5 утра каждый день
0 5 * * * /backup.sh >> /var/log/cron.log 2>&1
